# pacefinder.service
[Unit]
Description=Pacefinder UDP Listener
After=network.target local-fs.target
# Listener falls back to local storage if USB isn't mounted, so no hard USB dep.

# Manual copy-paste template — assumes you cloned to /home/pi/pacefinderapp
# and run as user "pi". If either differs, edit WorkingDirectory/ExecStart/User
# below, OR just run ./install-pi.sh instead, which stamps these automatically
# and verifies the service comes up.
[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pacefinderapp
ExecStart=/usr/bin/python3 /home/pi/pacefinderapp/listener.py
Restart=always
RestartSec=5
# Never stop restarting regardless of how many failures occur.
StartLimitIntervalSec=0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
